# apache/cbo-flask.service
[Unit]
Description=Central Bank of Oman Chatbot Flask Backend
Documentation=https://github.com/your-repo/cbo-chatbot
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cbo-chatbot-backend
ExecStart=/usr/bin/python3 run_app.py
Restart=always
RestartSec=10
Environment=FLASK_ENV=production
Environment=FLASK_SECRET_KEY=your-production-secret-key
Environment=DB_NAME=doc_analyzer

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cbo-flask

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/cbo-chatbot-backend

[Install]
WantedBy=multi-user.target
